# systemd/console-1706.service
[Unit]
Description=console-1706 local web UI
After=network.target

[Service]
WorkingDirectory=%h/projects/console-1706
ExecStart=%h/projects/console-1706/.venv/bin/console-1706 serve --config %h/.config/console-1706/config.yml
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
