# systemd/console-1701-scan.timer
[Unit]
Description=Run console-1701 scanner every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true
Unit=console-1701-scan.service

[Install]
WantedBy=timers.target
