# systemd/console-1701.service
[Unit]
Description=console-1701 local web UI
After=network.target

[Service]
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/console-1701 serve --config %h/.config/console-1701/config.yml
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
